# roles/k8s_node_setup/templates/udev-trigger.service.j2
[Unit]
Description=Run udevadm trigger every 10 seconds

[Service]
Restart=always
ExecStart=/bin/bash -c "while true; do udevadm trigger ; sleep 10; done"
